# backend/systemd/worldmap.timer
[Unit]
Description=update map all 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
